# scripts/service/vm-net-init-http/run
#!/bin/sh

set -eu

cd /etc/vm-net-init/
exec 2>&1
exec /usr/bin/python3 -m http.server 1337
